# deploy/server/cdc-realtime-streaming.service
[Unit]
Description=CDC Warehouse Spark Structured Streaming Kafka To Kudu
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cdc
Group=cdc
WorkingDirectory=/opt/cdc-warehouse-platform
EnvironmentFile=/etc/cdc-warehouse/jobs.env
ExecStart=/opt/cdc-warehouse-platform/deploy/server/run_realtime_streaming.sh
Restart=always
RestartSec=10
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
