# scripts/com.cdcwarehouse.ops-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cdcwarehouse.ops-refresh</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/liuming/Documents/Codex/2026-07-07/wo/cdc-warehouse-platform/scripts/refresh_ops_snapshot.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/liuming/Documents/Codex/2026-07-07/wo/cdc-warehouse-platform</string>
  <key>StartInterval</key>
  <integer>5</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/liuming/Documents/Codex/2026-07-07/wo/cdc-warehouse-platform/data/ops/launchd.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/liuming/Documents/Codex/2026-07-07/wo/cdc-warehouse-platform/data/ops/launchd.err</string>
</dict>
</plist>
